# etc/cron.d/eis
# vi: ft=crontab :
MAILTO=root
*/5 * * * * root /usr/local/bin/eis-import-all.sh
